# docker/supervisord.allinone.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:postgres]
command=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/data
user=postgres
autorestart=true

[program:redis]
command=redis-server --protected-mode no
autorestart=true

[program:asterisk]
command=/usr/sbin/asterisk -f
autorestart=true

[program:kamailio]
command=/usr/sbin/kamailio -DD -E
autorestart=true

[program:api]
directory=/opt/api
command=node app.js
autorestart=true
environment=PORT="3000"

[program:dashboard]
directory=/opt/dashboard
command=npm run start
autorestart=true
environment=PORT="3001"
